# pibot-reminder.service
[Unit]
Description=AI Reminder System on Raspberry Pi
After=graphical-session.target sound.target

[Service]
Type=simple
WorkingDirectory=/home/eric/pi_aireminder
Environment="DISPLAY=:0"
Environment="XAUTHORITY=/home/eric/.Xauthority"
ExecStart=/home/eric/pi_aireminder/venv/bin/python /home/eric/pi_aireminder/src/main.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
